# .config/autostart/flameshot.desktop
[Desktop Entry]
Type=Application
NoDisplay=true
Terminal=false
Name=Flameshot
Exec=/usr/bin/flameshot
Comment=Screenshot Utility
Categories=System;
